# crontab
# Update every minute
* * * * * /root/cron_task.sh /root/update_haproxy.sh

